# tabs/helium.cron
MAILTO=pheenan@ribometrix.com
# ASMS database backup, retired on 2021-07-12
#59 22 * * *  cd /home/pheenan/src_prh/Ribometrix/PrhCompoundManager/Code/mains/backup/ ; /usr/bin/bash pg_backup_rotated.sh -c pg_backup.config
# ASMS QC tracking, retired on 2021-07-12
#01 6,15,22 * * * cd /home/pheenan/mnt/helium/arch/mystuff/src_prh/PeakSumsAndHits/Code/control_tracking ; /home/pheenan/anaconda3/bin/python main_track_controls.py --out_directory /home/pheenan/mnt/helium/arch/rndgrp/ASMS/QC_tracking/ 
# ASMS database updating, retired on 2021-07-12
#59 3 * * * cd /home/pheenan/mnt/helium/arch/mystuff/src_prh/PeakSumsAndHits/Code/database/mains/ ; /home/pheenan/anaconda3/bin/python _main_adaptively_add_experiments.py --in_directory ../../../Experiments --forgiving TRUE --overwrite TRUE
# compound manager, every 6 hours, monday to friday
59 22,4,10,16 * * 1-5 /home/pheenan/mnt/helium/arch/mystuff/src_prh/PrhCron/sh/cronic /home/pheenan/src_prh/Ribometrix/PrhCompoundManager/Code/mains/main_of_mains/main_helium_production.sh
# 1 AM on saturdays: re-read everything
# 0 1 * *  6 /usr/bin/bash /home/pheenan/src_prh/Ribometrix/PrhCompoundManager/Code/mains/_main_reread_regmol/main.sh